# systemd/llama-gemma4-e4b.service
[Unit]
Description=llama.cpp server — Gemma 4 E4B (coordinator/general)
After=default.target

[Service]
# Vulkan (RADV) backend, NOT ROCm: on gfx1151 the ROCm/HIP runtime keeps the
# GPU spinning at 100% busy / max clock / ~35W even with zero requests, 24/7.
# Prime suspect in the recurring silent hangs — see
# docs/incidents/INVESTIGATION_2026-06-09_hangs.md. Vulkan idles correctly
# (0% / ~7W) at a ~15-20% tok/s cost. Switched 2026-06-09.
Environment="LD_LIBRARY_PATH=/home/drew/pai/pai_workspace/llama-cpp/llama-b9585-vulkan"
ExecStart=/home/drew/pai/pai_workspace/llama-cpp/llama-b9585-vulkan/llama-server \
    -m /opt/models/google/gemma-4-E4B-it-qat-UD-Q4_K_XL.gguf \
    --host 127.0.0.1 \
    --port 11438 \
    -ngl 99 \
    --flash-attn on \
    --ctx-size 32768 \
    --cache-reuse 256 \
    --swa-full \
    --reasoning-budget 256
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
