# systemd/llama-devstral.service
[Unit]
Description=llama.cpp server — Devstral-2512 (coding/devops)
After=default.target

[Service]
Environment="HSA_OVERRIDE_GFX_VERSION=11.5.1"
Environment="LD_LIBRARY_PATH=/usr/local/lib/ollama/rocm:/home/drew/pai/pai_workspace/llama-cpp/llama-gfx1151-rocm7"
ExecStart=/home/drew/pai/pai_workspace/llama-cpp/llama-gfx1151-rocm7/llama-server \
    -m /opt/models/mistralai/mistralai_Devstral-Small-2-24B-Instruct-2512-Q8_0.gguf \
    --host 127.0.0.1 \
    --port 11436 \
    -ngl 99 \
    --flash-attn on \
    --ctx-size 32768 \
    --cache-reuse 256 \
    --chat-template-file /opt/models/mistralai/devstral-template-permissive.jinja
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
