# systemd/llama-gemma4-e4b.service
[Unit]
Description=llama.cpp server — Gemma 4 E4B (coordinator/general)
After=default.target

[Service]
Environment="HSA_OVERRIDE_GFX_VERSION=11.5.1"
Environment="LD_LIBRARY_PATH=/home/drew/pai/pai_workspace/llama-cpp/llama-b8746:/usr/local/lib/ollama/rocm"
ExecStart=/home/drew/pai/pai_workspace/llama-cpp/llama-b8746/llama-server \
    -m /opt/models/google/google_gemma-4-E4B-it-Q4_K_M.gguf \
    --host 127.0.0.1 \
    --port 11438 \
    -ngl 99 \
    --flash-attn on \
    --ctx-size 32768 \
    --cache-reuse 256
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
